# leadflow.service
[Unit]
Description=LeadFlow Django ASGI Application
After=network.target postgresql.service

[Service]
Type=simple
# Replace 'your-plesk-user' with your actual Plesk user
User=your-plesk-user
Group=psacln

# Working directory
WorkingDirectory=/var/www/vhosts/your-domain.com/httpdocs/backend

# Environment
Environment="PATH=/opt/plesk/python/3.12/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DJANGO_SETTINGS_MODULE=backend.settings"

# Load environment variables from file (if you create one)
# EnvironmentFile=/var/www/vhosts/your-domain.com/httpdocs/.env

# Command to start the application
# Adjust Python path based on your Plesk Python installation
ExecStart=/opt/plesk/python/3.12/bin/daphne -b 127.0.0.1 -p 8000 backend.asgi:application

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=leadflow

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

